# debian/libapache2-mod-authnz-restauth.postinst
#!/bin/sh

set -e

case "$1" in
    configure)
        echo -n "Removing old restauth module name authnz-restauth.load.."
		rm -f /etc/apache2/mods-available/authnz-restauth.load
        invoke-rc.d apache2 reload
        ;;
esac

#DEBHELPER#

